# mysql/supervisord.conf
[supervisord]
nodemon=true

[program:mysqld]
command=mysqld
stdout_lgofile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0