# mysql/supervisord.conf
[supervisord]
nodemon=true
loglevel=debug

[program:mysql]
command=/usr/bin/mysqld_safe
user=mysql
priority=0